# bootsy-cmd/src/main/resources/kube-proxy.service
[Unit]
Description=kube-proxy

[Service]
ExecStart=/usr/bin/kube-proxy --bind-address=$bind_address --master=$api_server_endpoint --hostname-override=$ip_address --cluster-cidr=$cluster_cidr --kubeconfig=$kubeconfig
Restart=always

[Install]
WantedBy=multi-user.target
